# systemd_files/update_snapshot_entry.timer
[Unit]
Description=Run update_snapshot_entry after a good boot

[Timer]
OnBootSec=2
Persistent=true

[Install]
WantedBy=timers.target
